# roles/docker_nginx/templates/etc_systemd_system_certbot.timer.j2
[Unit]
Description=Run certbot daily

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
